# deploy/oh-my-cassette-web.service.example
[Unit]
Description=Oh My Cassette Web Demo
After=network-online.target
Wants=network-online.target

# The browser UI must be built before this unit starts (build time only, needs
# Node.js + npm). Run during provisioning / on each deploy:
#   /home/ubuntu/oh-my-cassette/web_demo/build_frontend.sh
# This produces web_demo/frontend/dist, which the server serves under /static.
[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/oh-my-cassette
EnvironmentFile=-/etc/oh-my-cassette-web.env
Environment=HOME=/home/ubuntu
Environment=XDG_CACHE_HOME=/home/ubuntu/.cache
Environment=PLAYWRIGHT_BROWSERS_PATH=/home/ubuntu/.cache/ms-playwright
ExecStart=/home/ubuntu/oh-my-cassette/.venv-web/bin/python -m web_demo.server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
